# php/supervisor/cron.conf
[program:cron]
autostart = true
autorestart = true
command=crond -f
stderr_logfile=/var/www/log/cron.err.log
stdout_logfile=/var/www/log/cron.out.log
priority=300
