# php/supervisor/php.conf
[program:php]
command=docker-php-entrypoint php-fpm
stderr_logfile=/var/www/log/php-fpm_consumer.err.log
stdout_logfile=/var/www/log/php-fpm_consumer.out.log
